# supervisor/predestination.conf
[program:predestination]
command=/home/ubuntu/predestination/web
environment=PORT=8080
user=web
