# supervisor/predestination.conf
[program:predestination]
command=/var/www/predestination/web
environment=PORT="8080",TRANSPORTS="polling"
user=web
